# deploy/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=/usr/local/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    DATABASE_URL="sqlite+aiosqlite:////data/db/documind.db",
    UPLOAD_DIR="/data/uploads",
    QDRANT_MODE="local",
    QDRANT_LOCAL_PATH="/data/qdrant_storage",
    HOME="/app"

[program:frontend]
command=/usr/bin/npx next start -p 3000
directory=/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    NEXT_PUBLIC_API_URL="",
    NEXT_PUBLIC_WS_URL="",
    NODE_ENV="production"
